# ansible/roles/fdg.gunicorn/templates/gunicorn.socket.j2
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn/socket
User=www-data

[Install]
WantedBy=sockets.target
